# node_exporter.ini
[program:node_exporter]
command=/opt/node_exporter/node_exporter
user=exporter
directory=/opt/node_exporter/node_exporter
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stdout_logfile_backups=0
stdout_logfile_maxbytes=50MB
